# invidious/rootfs/etc/services.d/invidious/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Cleanup when Invidious service stops
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

declare APP_EXIT_CODE=${1}

bashio::log.info "Invidious service stopped with exit code ${APP_EXIT_CODE}"

# Stop PostgreSQL gracefully
if s6-setuidgid postgres pg_ctl -D /data/postgresql status > /dev/null 2>&1; then
    bashio::log.info "Stopping PostgreSQL..."
    s6-setuidgid postgres pg_ctl -D /data/postgresql -w stop
fi

# Handle service exit codes
# Exit code 256 (128+128) typically indicates normal s6 service termination
# Exit code 0 indicates successful completion
# Any other exit code indicates an error that should halt the service
if [[ "${APP_EXIT_CODE}" -ne 0 ]] && [[ "${APP_EXIT_CODE}" -ne 256 ]]; then
  bashio::log.error "Invidious exited with error code ${APP_EXIT_CODE}"
  bashio::log.error "Check the logs above for details on what went wrong"
  bashio::log.error "Halting service to prevent restart loop"
  # Exit with 125 to signal s6-overlay to halt the service and not restart it
  exit 125
fi
