# invidious/rootfs/etc/services.d/invidious/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Invidious service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Starting Invidious add-on..."

# Initialize PostgreSQL if not already done
if [ ! -d "/data/postgresql" ]; then
    bashio::log.info "Initializing PostgreSQL database..."
    mkdir -p /data/postgresql
    chown postgres:postgres /data/postgresql
    su-exec postgres initdb -D /data/postgresql
    
    # Start PostgreSQL temporarily to create database
    su-exec postgres pg_ctl -D /data/postgresql -o "-c listen_addresses='localhost'" -w start
    
    # Create Invidious database and user
    # Note: Password is hardcoded since PostgreSQL only listens on localhost
    # and is only accessible within the container. No external access is possible.
    su-exec postgres psql -c "CREATE USER invidious WITH PASSWORD 'invidious';"
    su-exec postgres psql -c "CREATE DATABASE invidious OWNER invidious;"
    
    # Stop PostgreSQL
    su-exec postgres pg_ctl -D /data/postgresql -w stop
    
    bashio::log.info "PostgreSQL initialization complete"
fi

# Start PostgreSQL in the background
bashio::log.info "Starting PostgreSQL..."
su-exec postgres pg_ctl -D /data/postgresql -o "-c listen_addresses='localhost'" -w start

# Generate HMAC key if not exists
if [ ! -f "/data/hmac_key" ]; then
    bashio::log.info "Generating HMAC key..."
    pwgen -s 32 1 > /data/hmac_key
fi

HMAC_KEY=$(cat /data/hmac_key)

# Create Invidious configuration
bashio::log.info "Creating Invidious configuration..."
cat > /invidious/config/config.yml << EOF
database_url: postgres://invidious:invidious@localhost:5432/invidious
check_tables: true
port: 3000
host_binding: 0.0.0.0
hmac_key: "${HMAC_KEY}"
EOF

# Wait for PostgreSQL to be ready
bashio::log.info "Waiting for PostgreSQL to be ready..."
for i in {1..30}; do
    if su-exec postgres pg_isready -q; then
        bashio::log.info "PostgreSQL is ready"
        break
    fi
    if [ $i -eq 30 ]; then
        bashio::log.error "PostgreSQL failed to start"
        exit 1
    fi
    sleep 1
done

# Run Invidious
bashio::log.info "Starting Invidious on port 3000..."
cd /invidious
exec /invidious/invidious
