# invidious/rootfs/etc/services.d/invidious/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Invidious service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Starting Invidious add-on..."

# Initialize PostgreSQL if not already done
if [ ! -d "/data/postgresql" ]; then
    bashio::log.info "Initializing PostgreSQL database..."
    mkdir -p /data/postgresql || {
        bashio::log.error "Failed to create PostgreSQL data directory"
        exit 1
    }
    chown -R postgres:postgres /data/postgresql
    # Use s6-setuidgid instead of su-exec
    s6-setuidgid postgres initdb -D /data/postgresql || {
        bashio::log.error "Failed to initialize PostgreSQL database"
        exit 1
    }
    
    # Start PostgreSQL temporarily to create database
    s6-setuidgid postgres pg_ctl -D /data/postgresql -o "-c listen_addresses='localhost'" -w start || {
        bashio::log.error "Failed to start PostgreSQL for initialization"
        exit 1
    }
    
    # Create Invidious database and user
    # Note: Password is hardcoded since PostgreSQL only listens on localhost
    # and is only accessible within the container. No external access is possible.
    s6-setuidgid postgres psql -c "CREATE USER invidious WITH PASSWORD 'invidious';" || {
        bashio::log.error "Failed to create Invidious database user"
        s6-setuidgid postgres pg_ctl -D /data/postgresql -w stop
        exit 1
    }
    s6-setuidgid postgres psql -c "CREATE DATABASE invidious OWNER invidious;" || {
        bashio::log.error "Failed to create Invidious database"
        s6-setuidgid postgres pg_ctl -D /data/postgresql -w stop
        exit 1
    }
    
    # Stop PostgreSQL
    s6-setuidgid postgres pg_ctl -D /data/postgresql -w stop
    
    bashio::log.info "PostgreSQL initialization complete"
fi

# Start PostgreSQL in the background
bashio::log.info "Starting PostgreSQL..."
s6-setuidgid postgres pg_ctl -D /data/postgresql -o "-c listen_addresses='localhost'" -w start || {
    bashio::log.error "Failed to start PostgreSQL"
    exit 1
}

# Generate HMAC key if not exists
if [ ! -f "/data/hmac_key" ]; then
    bashio::log.info "Generating HMAC key..."
    pwgen -s 32 1 > /data/hmac_key || {
        bashio::log.error "Failed to generate HMAC key"
        exit 1
    }
fi

HMAC_KEY=$(cat /data/hmac_key) || {
    bashio::log.error "Failed to read HMAC key"
    exit 1
}

if [ -z "$HMAC_KEY" ]; then
    bashio::log.error "HMAC key is empty"
    exit 1
fi

# Create Invidious configuration
bashio::log.info "Creating Invidious configuration..."
mkdir -p /invidious/config || {
    bashio::log.error "Failed to create Invidious config directory"
    exit 1
}
cat > /invidious/config/config.yml << EOF
database_url: postgres://invidious:invidious@localhost:5432/invidious
check_tables: true
port: 3000
host_binding: 0.0.0.0
hmac_key: "${HMAC_KEY}"
EOF

# Wait for PostgreSQL to be ready
bashio::log.info "Waiting for PostgreSQL to be ready..."
for i in {1..30}; do
    if s6-setuidgid postgres pg_isready -q; then
        bashio::log.info "PostgreSQL is ready"
        break
    fi
    if [ "$i" -eq 30 ]; then
        bashio::log.error "PostgreSQL failed to start"
        exit 1
    fi
    sleep 1
done

# Verify Invidious binary exists
if [ ! -f "/invidious/invidious" ]; then
    bashio::log.error "Invidious binary not found at /invidious/invidious"
    exit 1
fi

if [ ! -x "/invidious/invidious" ]; then
    bashio::log.error "Invidious binary is not executable"
    exit 1
fi

# Run Invidious
bashio::log.info "Starting Invidious on port 3000..."
cd /invidious || {
    bashio::log.error "Failed to change to Invidious directory"
    exit 1
}

# Log configuration for debugging
bashio::log.info "Invidious configuration:"
bashio::log.info "Database: postgres://invidious:***@localhost:5432/invidious"
bashio::log.info "Port: 3000"
bashio::log.info "Host: 0.0.0.0"

exec /invidious/invidious
